# ducc/unitfiles/cvmfs_ducc.service.in
[Unit]
Description=DUCC Daemon Unpacking Containers in CVMFS
After=network.target dockerd.service
Requires=dockerd.service

[Service]
Environment="DUCC_RECIPE_FILE=YourRecipeFile_UPDATE_THIS_VALUE.yaml"
Environment="DUCC_DOCKER_REGISTRY_PASS=ThePassOfTheUserForTheImageRegistry_UPDATE_THIS_VALUE"

ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cvmfs_ducc loop $RECIPE_FILE

[Install]
WantedBy=multi-user.target
